# scripts/vless-manager.service.example
[Unit]
Description=VLESS Proxy Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vless-proxy
EnvironmentFile=/root/vless-proxy/.env
ExecStart=/root/vless-proxy/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-manager

[Install]
WantedBy=multi-user.target
