# scripts/vless-manager.service
[Unit]
Description=VLESS Proxy Manager
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/vless-proxy
EnvironmentFile=/path/to/vless-proxy/.env
ExecStart=/path/to/vless-proxy/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-manager

[Install]
WantedBy=multi-user.target
